# systemd/dmarcmailcheck.timer
[Unit]
Description=Check mail for new dmarc reports on schedule

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target