# systemd/dmarcmailcheck.service
[Unit]
Description=Check mail for new dmarc reports

[Service]
Type=simple
Nice=10
ExecStart=php {PATH}/bin/console app:checkmailbox