# dockerfiles/etc/periodic/daily/checkmail.sh
#!/bin/sh
php /var/www/html/bin/console app:checkmailbox